# deploy/systemd/pdf-template-editor.service
[Unit]
Description=PDF Template Editor
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/var/www/pdf-template-editor
EnvironmentFile=/var/www/pdf-template-editor/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

# Základní hardening. Pokud by některé omezení vadilo zápisu/cache, upravit podle serveru.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/pdf-template-editor

[Install]
WantedBy=multi-user.target
